# telescope/assets/services/webhook.socket
[Unit]
Description=Webhook server socket

[Socket]
ListenStream=51337

[Install]
WantedBy=multi-user.target
