# telescope/assets/services/webhook.socket
[Unit]
Description=Webhook server socket

[Socket]
BindIPv6Only=ipv6-only
ListenStream=51337

[Install]
WantedBy=multi-user.target
